# data/systemd/rnaseq-pipeline-cleanup.service
[Unit]
Description=PavLab RNA-Seq Pipeline Cleanup

[Service]
Type=simple
ExecStart=/bin/bash -c "./scripts/remove-old-data"
