# data/systemd/rnaseq-pipeline-cleanup.timer
[Unit]
Description=Cleanup the RNA-Seq Pipeline hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
